# templates/etc/init.d/mongodb-mms-automation.erb
#!/bin/bash
#
# mongodb-mms This starts and stops mongodb-mms
### BEGIN INIT INFO
# Provides: mongodb-mms
# Required-Start: $network $local_fs $remote_fs
# Required-Stop: $network $local_fs $remote_fs
# Should-Start: $named
# Should-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: An agent to raise mongodb stats in MMS
# Description: MongoDB Monitoring Service
### END INIT INFO
###
### chkconfig: 2345 12 88
### description: mongodb-mms is mongodb-mms
### processname:mongodb-mms
### pidfile: /var/run/mongodb-mms.pid
##### BEGIN INIT INFO
### Provides: $mongodb-mms
##### END INIT INFO
### Source function library.
binary="/opt/mms/mongodb-mms.pl"
[ -x $binary ] || exit 0
RETVAL=0
start() {
    echo -n "Starting mongodb-mms: "
    $binary start
    RETVAL=$?
    PID=$!
    echo
    [ $RETVAL -eq 0 ] && touch /var/mms
    echo $PID > /var/run/mongodb-mms.pid
}
stop() {
    echo -n "Shutting down mongodb-mms: "
    $binary stop
    killall mongodb-mms
    RETVAL=$?
    echo
    if [ $RETVAL -eq 0 ]; then
        rm -f /var/lock/subsys/mongodb-mms
        rm -f /var/run/mongodb-mms.pid
    fi
}
restart() {
    echo -n "Restarting mongodb-mms: "
    $binary stop
    sleep 2
    $binary start
}
case "$1" in
    start)
        start
    ;;
    stop)
        stop
    ;;
    status)
        $binary status
    ;;
    restart)
        restart
    ;;
    *)
        echo "Usage: $0 {start|stop|status|restart}"
    ;;
esac
exit 0
